# build/package/debian/lib/systemd/system/pbs-d2d-backup.service
[Unit]
Description=Proxmox D2D Backup Proxy
Wants=network-online.target
After=network.target
Wants=proxmox-backup-proxy.service
After=proxmox-backup-proxy.service

[Service]
Type=simple
ExecStart=/usr/bin/pbs-d2d-backup
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/proxmox-backup/pbs-d2d.pid
Restart=on-failure
User=root
Group=root

[Install]
WantedBy=multi-user.target

