# build/package/debian/DEBIAN/postinst
#!/bin/bash
set -e

chmod +x /usr/bin/pbs-plus

systemctl daemon-reload

# Enable the service to start on boot
systemctl enable pbs-plus.service

# Start the service
systemctl start pbs-plus.service

exit 0
