# build/package/linux_agent/debian/lib/systemd/system/pbs-plus-agent.service
[Unit]
Description=PBS Plus Agent
Wants=network-online.target
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/pbs-plus-agent
ExecReload=/bin/kill -HUP $MAINPID
PIDFile=/run/proxmox-backup/pbs-plus-agent.pid
Restart=on-failure
User=root
Group=root

[Install]
WantedBy=multi-user.target

