# build/debian/DEBIAN/postinst
#!/bin/bash
set -e

chmod +x /usr/bin/pbs-d2d-backup

systemctl daemon-reload

# Enable the service to start on boot
systemctl enable pbs-d2d-backup.service

# Start the service
systemctl start pbs-d2d-backup.service

exit 0
